# systemd/morning-brief.service.template
[Unit]
Description=Morning Brief — generate and display today's digest
# Wait for the graphical session so xdg-open / notify-send work
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
# __PROJECT_DIR__ is replaced by install.sh with the real path at install time.
ExecStart=__PROJECT_DIR__/.venv/bin/python __PROJECT_DIR__/brief.py --wait-for-ollama
WorkingDirectory=__PROJECT_DIR__
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
